# bigtop-packages/src/common/hadoop/hadoop-hdfs-zkfc.service
[Unit]
Description=The Hadoop HDFS zkfc daemon
After=network.target
After=NetworkManager.target

[Service]
Type=simple
ExecStart=/usr/lib/hadoop/libexec/hadoop-hdfs-datanode start
ExecStop=/usr/lib/hadoop/libexec/hadoop-hdfs-datanode stop
PIDFile=/var/run/hadoop-hdfs/hadoop-hdfs-zkfc.pid
SuccessExitStatus=143

#######################################
# Note: Below are cgroup options
#######################################
#Delegate=true
#CPUAccounting=true
#CPUShares=1024

[Install]
WantedBy=multi-user.target
